# profiles/sovereign-ai.target
[Unit]
Description=FabricTreaty Sovereign AI Profile (opt-in)
Wants=fabricd.service
After=multi-user.target
